# pulse-frontend.service
[Unit]
Description=Pulse Frontend (Vite)
After=network.target

[Service]
Type=simple
User=pulse
WorkingDirectory=/opt/pulse/frontend-modern
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=3
StandardOutput=append:/opt/pulse/frontend.log
StandardError=append:/opt/pulse/frontend.log
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target